# roles/prometheus/templates/node_exporter.service.j2
# Node Exporter Systemd Service
# Generated by Ansible - {{ ansible_date_time.iso8601 }}
# Author: hyhoncute

[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_group }}

ExecStart=/usr/local/bin/node_exporter \
  --web.listen-address={{ node_exporter_listen_address }}

SyslogIdentifier=node_exporter
Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/proc /sys

[Install]
WantedBy=multi-user.target
