# roles/prometheus/templates/prometheus.service.j2
# Prometheus Systemd Service
# Generated by Ansible - {{ ansible_date_time.iso8601 }}
# Author: hyhoncute

[Unit]
Description=Prometheus Monitoring System
Documentation=https://prometheus.io/docs/introduction/overview/
After=network.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_group }}

ExecStart={{ prometheus_install_dir }}/prometheus \
  --config.file={{ prometheus_config_dir }}/prometheus.yml \
  --storage.tsdb.path={{ prometheus_data_dir }} \
  --storage.tsdb.retention.time={{ prometheus_retention_time }} \
  {% if prometheus_retention_size != '0' %}
  --storage.tsdb.retention.size={{ prometheus_retention_size }} \
  {% endif %}
  --web.console.templates={{ prometheus_config_dir }}/consoles \
  --web.console.libraries={{ prometheus_config_dir }}/console_libraries \
  --web.listen-address={{ prometheus_listen_address }}

Restart=always
RestartSec=5

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ prometheus_data_dir }}

[Install]
WantedBy=multi-user.target